# init/home-assistant.service
[Unit]
Description=Home Assistant
After=multi-user.target alsa-restore.service 


[Service]
Type=simple
User=jambula
ExecStart=/bin/bash -c "/usr/local/bin/python3 /usr/local/bin/hass -c /etc/JambulaTV/homeassistant"

[Install]
WantedBy=graphical.target
