# init/jambulatv-tvguide-fta-generic.timer
[Unit]
Description=generation of daily EPG (Generic) data for free-to-air TV


[Timer]
OnCalendar=daily


[Install]
WantedBy=timers.target
