# init/reminders.timer
[Unit]
Description = 2 Hours Timer

[Timer]
OnBootSec = 5min
OnUnitActiveSec = 120min

[Install]
WantedBy = multi-user.target
